# systemd/mam-bonus-manager.timer
[Unit]
Description=Run MAM bonus manager periodically

[Timer]
# First run 5 minutes after boot, then every hour.
# The script only buys wedges when WEDGE_HOURS have elapsed.
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
Unit=mam-bonus-manager.service

[Install]
WantedBy=timers.target
